# configs/systemd/simple-sub.service
[Unit]
Description=VPN Subscription Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-subscription
ExecStart=/opt/vpn-subscription/venv/bin/python /opt/vpn-subscription/simple-sub.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target